# systemd/handex.service
[Unit]
Description=Handex Human-in-the-Loop Workspace
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/handex
EnvironmentFile=/etc/handex/handex.env
ExecStart=/opt/handex/.venv/bin/python -m handex.server
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
